# checker/crontab
SHELL=/bin/bash
*/5 * * * * $HOME/checker
